# launchd/com.openclaw.pocket-portal-health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openclaw.pocket-portal-health</string>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <!-- Edit URL + label to match your config -->
    <string>curl -fsS --max-time 2 http://127.0.0.1:4377/api/health &gt;/dev/null || launchctl kickstart -k gui/$UID/com.openclaw.pocket-portal</string>
  </array>

  <key>StandardOutPath</key>
  <string>/tmp/pocket-portal-health.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/pocket-portal-health.err.log</string>
</dict>
</plist>
